# systemd/instant-replay-kiosk.service
[Unit]
Description=Instant Replay touch UI (Chromium kiosk)
After=graphical.target replay-engine.service
Wants=replay-engine.service

[Service]
Type=simple
Environment=DISPLAY=:0
ExecStart=/usr/bin/chromium --kiosk --noerrdialogs --disable-infobars --app=http://127.0.0.1:8080/
Restart=on-failure
RestartSec=5

[Install]
WantedBy=graphical.target
