# systemd/replay-engine.service
[Unit]
Description=Instant Replay Engine
After=network.target graphical.target
Wants=network-online.target graphical.target
StartLimitIntervalSec=0
StartLimitBurst=0

[Service]
Type=simple
ExecStart=/opt/instant-replay/bin/replay-engine --appliance
Environment=DISPLAY=:0
Environment=XDG_SESSION_TYPE=x11
# XDG_RUNTIME_DIR set in replay-engine.service.d/override.conf (install-on-pi.sh)
Environment=GST_PLUGIN_PATH=/usr/lib/aarch64-linux-gnu/gstreamer-1.0
Restart=always
RestartSec=3
TimeoutStopSec=15
EnvironmentFile=-/etc/default/replay-engine
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=graphical.target
