# systemd/replay-engine.service
[Unit]
Description=Instant Replay Engine
After=network.target graphical.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/opt/instant-replay/bin/replay-engine --appliance
Environment=GST_PLUGIN_PATH=/usr/lib/aarch64-linux-gnu/gstreamer-1.0
Restart=always
RestartSec=3
StartLimitIntervalSec=0
EnvironmentFile=-/etc/default/replay-engine
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
